# src/DNSPod.sh
#!/bin/bash
#chkconfig: 2345 80 90
#description:auto_run
/service/DNSPod/DNSPod
